# roles/pull_timer/templates/harbor-yard.timer.j2
[Unit]
Description=Run Harbor Yard converge periodically

[Timer]
OnBootSec={{ harbor_yard_timer_onboot }}
OnUnitActiveSec={{ harbor_yard_timer_interval }}
RandomizedDelaySec={{ harbor_yard_timer_random_delay }}
Persistent=true

[Install]
WantedBy=timers.target
